# cloud-images/build/static_artifacts/master-node/bootstrap-lattice.service
[Unit]
# Wait until the docker is running and the kube-api is responsive
Requires=docker.service wait-for-kube-api.service
After=docker.service wait-for-kube-api.service

[Service]
Type=oneshot
ExecStart=/opt/lattice/bin/master-bootstrap-lattice

[Install]
WantedBy=multi-user.target
