# cloud-images/build/static_artifacts/master-node/register-master-node-dns.service
[Unit]
Before=kubelet.service

# Wait until the network is up to try to register dns
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/lattice/bin/master-register-master-node-dns
