# cloud-images/build/static_artifacts/master-node/var-opt-etcd.mount
[Unit]
# Don't mount the fs until it's been formatted.
After=format-var-opt-etcd.service
Requires=format-var-opt-etcd.service

[Mount]
What=/dev/xvdf
Where=/var/opt/etcd
Type=ext4
TimeoutSec=5m
